# frost-branding-base/etc/local.d/0PS1.start
#!/bin/sh

sed -i 's/PS1/#PS1/' /etc/skel/.bash /etc/bash.bashrc

rm -f $0
