# frost-branding-base/etc/local.d/local.stop
#!/bin/sh

[ -x /etc/rc.local.stop ] && /etc/rc.local.stop
