# frost-branding-base/etc/rc.local
#!/bin/sh

# /etc/rc.local -- rc.local for Frost Linux
#
# Enter your custom commands here.

